# ha-energy-optimizer/rootfs/etc/s6-overlay/s6-rc.d/energy-optimizer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e
bashio::log.info "Starting HA Energy Optimizer..."
bashio::config > /data/options.json
bashio::config.has_value "database.host" || { bashio::log.fatal "Database host not configured"; exit 1; }
bashio::config.has_value "homeassistant.token" || { bashio::log.fatal "HA token not configured"; exit 1; }
bashio::log.info "Configuration OK"
cd /app
exec python3 main.py
