# ha-energy-optimizer/rootfs/etc/s6-overlay/s6-rc.d/energy-optimizer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# Service run script for s6-overlay.
# Service run-script voor s6-overlay.
set -e

bashio::log.info "Starting HA Energy Optimizer / HA Energy Optimizer wordt gestart..."

# Write HA options to /data/options.json
# Schrijf HA-opties naar /data/options.json
bashio::config > /data/options.json

# Validate required settings / Valideer vereiste instellingen
if ! bashio::config.has_value "database.host"; then
    bashio::log.fatal "Database host is not configured / Database host is niet ingesteld"
    exit 1
fi

if ! bashio::config.has_value "homeassistant.token"; then
    bashio::log.fatal "HA token is not configured / HA token is niet ingesteld"
    exit 1
fi

bashio::log.info "Configuration OK / Configuratie OK"

# Start the optimizer (foreground — s6 manages the process)
# Start de optimizer (voorgrond — s6 beheert het proces)
cd /app
exec python3 main.py
