# ha-energy-optimizer/rootfs/etc/services.d/energy-optimizer/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting HA Energy Optimizer..."

# Write HA configuration to options.json
# Schrijf HA-configuratie naar options.json
bashio::config > /data/options.json

bashio::log.info "Starting configuration GUI on port 8099..."
python3 -m flask --app /app/gui/app.py run \
    --host 0.0.0.0 \
    --port 8099 \
    --no-debugger &

bashio::log.info "Starting optimizer engine..."
exec python3 /app/main.py
