# ha-energy-optimizer/rootfs/etc/services.d/energy-optimizer/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting HA Energy Optimizer v0.2.0..."

OPTIONS_FILE="/data/options.json"

if [ -f "${OPTIONS_FILE}" ] && [ -s "${OPTIONS_FILE}" ]; then
    bashio::log.info "Using existing options.json"
else
    bashio::log.info "Fetching configuration from supervisor..."
    curl -s \
        -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
        http://supervisor/addons/self/options/config \
        > "${OPTIONS_FILE}" || echo '{}' > "${OPTIONS_FILE}"
fi

bashio::log.info "Configuration ready"

# Kill any existing GUI process on port 8099 before starting
# Beëindig eventueel bestaand GUI-proces op poort 8099
fuser -k 8099/tcp 2>/dev/null || true

bashio::log.info "Starting configuration GUI on port 8099..."
cd /app
python3 -m flask --app gui/app.py run \
    --host 0.0.0.0 \
    --port 8099 \
    --no-debugger &

bashio::log.info "Starting optimizer engine..."
exec python3 main.py
