# ha-energy-optimizer/ha-dockerfix/rootfs/etc/s6-overlay/s6-rc.d/energy-optimizer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "Starting HA Energy Optimizer..."

bashio::config > /data/options.json

if ! bashio::config.has_value "database.host"; then
    bashio::log.fatal "Database host not configured"
    exit 1
fi

if ! bashio::config.has_value "homeassistant.token"; then
    bashio::log.fatal "HA token not configured"
    exit 1
fi

bashio::log.info "Configuration OK — starting optimizer engine..."
cd /app
exec python3 main.py
