# rootfs/etc/s6-overlay/s6-rc.d/gui-optimizer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# GUI service run script.
# GUI service run-script.
set -e

bashio::log.info "Starting configuration GUI on port 8099 / Configuratie-GUI starten op poort 8099..."

cd /app
exec python3 -m flask --app gui/app.py run \
    --host 0.0.0.0 \
    --port 8099 \
    --no-debugger
